# com.hey-claude.chotu.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Generated from this template by scripts/setup.sh (paths substituted). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hey-claude.chotu</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>chotu</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/hey-claude/daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/hey-claude/daemon.err.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
