# com.hey-claude.chotu.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Generated from this template by scripts/setup.sh (paths substituted).
     Runs an internal-disk launcher (launch-chotu.sh) that waits for the repo's
     venv before starting chotu, so auto-start survives login even when the repo
     + venv live on an external volume. No external path is referenced here. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hey-claude.chotu</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__LAUNCHER__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on a crash (non-zero exit), but NOT on a clean exit(0). chotu exits 0
         when the mic device disconnects (Bluetooth drop) so it stops instead of running
         deaf; SuccessfulExit=false leaves it down for a manual reconnect + restart. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/hey-claude/daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/hey-claude/daemon.err.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
